# recipes-common/meminsight/files/meminsight-runner.service
[Unit]
Description=MemInsight Memory Stats Runner

[Service]
Type=simple
ExecStart=/usr/bin/xMemInsight

[Install]
WantedBy=multi-user.target
