# recipes-thunder/wpe-framework/wpeframework/thunder.service.in
[Unit]
Description=Thunder RPC version-R5.3
Wants= local-fs.target
After= local-fs.target securemount.service
RequiresMountsFor=/opt/secure

[Service]
Type=notify
Environment="XDG_RUNTIME_DIR=/tmp"
Environment="GST_REGISTRY=/opt/.gstreamer/registry.bin"
Environment="GST_REGISTRY_UPDATE=no"
ExecStart=/usr/bin/Thunder -b
ExecStartPost=-/bin/sh -c '/bin/echo ${MAINPID} > /tmp/thunder.pid'
KillSignal=SIGKILL
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
