# recipes-extended/wpe-framework/wpeframework/wpeframework.service.in
[Unit]
Description=wpeframework
Wants=network-online.target local-fs.target
After= network-online.target local-fs.target

[Service]
Type=notify
EnvironmentFile=-/etc/default/westeros-env
Environment="LD_LIBRARY_PATH=/usr/lib:/media/apps/netflix/usr/lib:/tmp/netflix/usr/lib:/media/apps/libcobalt/usr/lib:/tmp/libcobalt/usr/lib:/usr/lib/javascriptcore"
Environment="XDG_RUNTIME_DIR=/tmp"
Environment="GST_REGISTRY=/opt/.gstreamer/registry.bin"
Environment="GST_REGISTRY_UPDATE=no"
Environment="ENABLE_WEBKITBROWSER_PLUGIN_ACCESSIBILITY=1"
Environment="TTS_USE_THUNDER_CLIENT=1"
Environment="XDG_DATA_HOME=/opt/QT/home"
Environment="FORCE_SVP=TRUE"
Environment="WESTEROS_GL_GRAPHICS_MAX_SIZE=1920x1080"
Environment="FORCE_SAP=TRUE"
ExecStart=-/usr/bin/WPEFramework -b
# Manually adding PID file support. This is used by Thunder clients to know if
# Thunder aka WPEFramework process is restarted. Please note if PIDFile name is changed
# it will need corresponding code change in ThunderClientLibraries (power_controller) too.
ExecStartPost=-/bin/sh -c '/bin/echo ${MAINPID} > /tmp/wpeframework.pid'
ExecStop=/usr/bin/killall -9 WPEFramework
ExecStopPost=/usr/bin/killall -9 WPEProcess
KillSignal=SIGKILL
Restart=always

[Install]
WantedBy=multi-user.target
